# s6-lamp/root/etc/services.d/mysql/run
#!/usr/bin/with-contenv sh

exec mysqld --user mysql --bind-address 0.0.0.0;
